# py/iqdjs/p-landlib/apps/xamppSsl/installer/data/x/systemctl/fastxampp.service
[Unit]
Description=Run fastxampp deamon

[Service]
Type=oneshot
ExecStart=/usr/local/fastxampp/fastxamppd.sh

[Install]
WantedBy=multi-user.target
